# deploy/tangra-asterisk.service
[Unit]
Description=Tangra Asterisk Service - Read-only CDR viewer for FreePBX
Documentation=https://github.com/go-tangra/go-tangra-asterisk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tangra-asterisk -c /etc/tangra-asterisk
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Environment file (create /etc/tangra-asterisk/env with your values)
EnvironmentFile=-/etc/tangra-asterisk/env

# Admin service registration
Environment=ADMIN_GRPC_ENDPOINT=
Environment=GRPC_ADVERTISE_ADDR=0.0.0.0:9800
Environment=MODULE_REGISTRATION_SECRET=

# MySQL DSNs to FreePBX (override defaults in /etc/tangra-asterisk/data.yaml).
# Required at runtime — without them the service exits on startup.
Environment=ASTERISK_CDR_DSN=
Environment=ASTERISK_CONFIG_DSN=

# Call recordings directory (FreePBX default: /var/spool/asterisk/monitor).
# When set, the module exposes GET /recordings/{linkedid} and the call
# drawer in the UI renders an inline audio player for any CDR row whose
# recordingfile column is non-empty. Leave empty to disable the endpoint.
Environment=ASTERISK_RECORDINGS_PATH=/var/spool/asterisk/monitor

# Prometheus base URL. When set, the Dashboards submenu in the UI is
# populated by proxying PromQL queries through this module to the
# configured Prometheus. Expects metrics from
# https://github.com/menta2k/freepbx-exporter under the `asterisk_*`
# namespace. Leave empty to disable the dashboards.
Environment=ASTERISK_PROMETHEUS_URL=http://127.0.0.1:9090

# mTLS server certs (gRPC server). Paths follow the cert-manager convention
# under $CERTS_DIR (set in /etc/tangra-asterisk/env):
#   {CERTS_DIR}/ca/ca.crt
#   {CERTS_DIR}/asterisk-server/server.{crt,key}
# Override with ASTERISK_CA_CERT_PATH / ASTERISK_SERVER_CERT_PATH /
# ASTERISK_SERVER_KEY_PATH in /etc/tangra-asterisk/env if needed.

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Allow read/write to config and certs only — module is read-only against
# FreePBX MySQL and embeds its own frontend, so no other writable paths.
ReadWritePaths=/etc/tangra-asterisk

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tangra-asterisk

[Install]
WantedBy=multi-user.target
